# abms-adminscripts/conf/backup/abms-backup.systemd.service
[Unit]
Description=ABMS backup service
After=pghelper.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

